# scripts/custom/monitor-startup.service
[Unit]
Description=Script

[Service]
ExecStart=/usr/bin/monitor-startup.sh

[Install]
WantedBy=sysinit.target
